# service/tlssnopper/run
#!/usr/bin/bash 
set -e
cd $( cd "$( dirname "${BASH_SOURCE[0]}" )" && pwd )
export DAEMONTOOLS_SERVICE=1
#:wsource ../../vars.sh
BCC_EXPERIMENTS_PATH=$(cd ../../ && pwd)
cd $BCC_EXPERIMENTS_PATH
cmd=$BCC_EXPERIMENTS_PATH/tlssnooper.sh
#source $BCC_EXPERIMENTS_PATH/.v/bin/activate
exec $cmd 2>/dev/null
