# service/exec_logger/run
#!/bin/bash
set -e
cd $( cd "$( dirname "${BASH_SOURCE[0]}" )" && pwd )
export DAEMONTOOLS_SERVICE=1
cmd=~root/bcc-experiments/execsnoop.sh
source ../../_service_utils.sh

cd / && eval $cmd 2>/dev/null \
	  | CREATE_LOG_FILTER "ssh" \
	  | CREATE_LOG_FILTER "python" \
	  | CREATE_LOG_FILTER "python2" \
	  | CREATE_LOG_FILTER "python3" \
	  | CREATE_LOG_FILTER "git" \


