# service/exec_logger/run
#!/bin/bash
set -e
cd $( cd "$( dirname "${BASH_SOURCE[0]}" )" && pwd )
export DAEMONTOOLS_SERVICE=1
cmd=~root/bcc-experiments/execsnoop.sh
source ../../_service_utils.sh
source ../../utils.sh

EXCLUDED_STRINGS="service/exec_logger|service/tcp_connection|bin/svstat service/|ls -alt service/"
REQUIRED_STRINGS='^{'

log_filter(){
  command egrep -v "$EXCLUDED_STRINGS" \
    | command grep "$REQUIRED_STRINGS"
}


pre_cmd="cd / && eval $cmd 2>/dev/null"


eval $(yaml2bash /root/bcc-experiments/etc/service_config.yaml)
log_filter_cmd="$(y2b_traverse Y2B|cut -d'=' -f2|xargs -I % echo -ne 'CREATE_LOG_FILTER %|' | sed 's/|$//g')"

cmd="$pre_cmd | log_filter | $log_filter_cmd"
eval $cmd


